# 01-linux-user-cron/cronjobs/daily_backup.cron
0 2 * * * /bin/bash /home/enoch/devops-projects-portfolio/01-linux-user-cron/scripts/backup_cronjob.sh >> /var/log/backup.log 2>&1
